# gomas.service
[Unit]
Description=Christmas lights

[Service]
Type=simple
User={{.User}}
ExecStart={{.BinaryPath}}
